# root/home/fyodor/.config/systemd/user/cloud-sync.timer
[Unit]
Description=Periodic cloud sync

[Timer]
OnStartupSec=5
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
