# root/home/fyodor/.config/systemd/user/cloud-sync.service
[Unit]
Description=Cloud sync (borg snapshot + rclone bisync)

[Service]
Type=oneshot
ExecStart=%h/.scripts/cloud-sync/sync --if-changed
